# udev/rules/99-cme-mkey.rules
ACTION=="add", ATTR{product}=="CME M-KEY", RUN+="/bin/su panther -c '/bin/bash ~/.scripts/udev/cmemkey.sh'"
