# udev/rules/99-vortex.rules
ACTION=="add", ENV{ID_TYPE}=="disk", ENV{ID_FS_LABEL}=="vortex", RUN+="/bin/mount /dev/%k /mnt/vortex"
