# applications/autostart/remmina-applet.desktop
[Desktop Entry]
Version=1.0
Name=Remmina Applet
Icon=remmina
Exec=env GDK_SCALE=1 remmina -i
Terminal=false
Type=Application
Hidden=false
RunHook=0
