# applications/autostart/rss_feeds.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=RSS feeds gatherer
Exec=.scripts/rss_feeds.sh
StartupNotify=false
Terminal=false
Hidden=false
